# templates/default/logstash-init.conf.erb
# This file is managed by Chef for <%= node['fqdn'] %>
# Do NOT modify this file directly.
#
# logstash - indexer instance
#
description     "logstash indexer instance"

start on virtual-filesystems
stop on runlevel [06]

respawn
respawn limit 5 30
limit nofile 65550 65550

env HOME=<%= @working_directory %>
env JAVA_OPTS='<%= @java_options %>'
chdir <%= @working_directory %>

# https://bugs.launchpad.net/upstart/+bug/812870
# We use `sudo -u ` instead…
#setuid <%= @user %>
#setgid <%= @group %>

script
exec sudo -u <%= @user %> <%= @binary_path %> <% if @web_interface %> web -- <% end %> agent -f <%= @configuration_directory %> --log <%= @log_directory %>/<%= @log_file %>
end script
